# PetaLinux/project-spec/meta-user/recipes-apps/chown-ebaz-uio-sdr-app/files/chown-ebaz-uio-sdr-app.service
[Unit]
Description=chown-ebaz-uio-sdr-app
  
[Service]
ExecStart=/usr/bin/chown-ebaz-uio-sdr-app
StandardOutput=journal+console
  
[Install]
WantedBy=multi-user.target
